# scripts/systemd/atlas-pool.service
[Unit]
Description=Atlas LVM thin pool
# The pool's PV is either real NVMe device(s) (Scaleway Elastic Metal) or a
# loopback PV (a sparse backing file on a stock droplet's root fs). A real-device
# PV survives a reboot intact; the loopback file survives but its loop binding
# does not, so the VG/pool are gone until the loop is re-attached. bootstrap-
# server.py creates the pool but is NOT re-run on boot, so this oneshot re-asserts
# it (ThinPool.ensure() → PoolBacking.reassert() re-binds the loop only when the
# backing is the file). Ordered before any VM unit (firecracker-vm@.service Wants
# network-online.target; this only needs local fs) so a VM never starts against a
# missing pool.
After=local-fs.target
Before=firecracker-vm@.service
DefaultDependencies=no
RequiresMountsFor=/var/lib/atlas/pool

[Service]
Type=oneshot
RemainAfterExit=yes
# Re-assert the pool via the durable atlas package. ThinPool().ensure() is
# idempotent: if the pool LV already exists it only re-attaches the loop device
# and activates the VG; otherwise it rebuilds pv/vg/pool from the backing file.
# The package lives under /var/lib/atlas/bin (placed by bootstrap), added to
# sys.path here so `from atlas.lvm import ThinPool` resolves.
# Runs under the Atlas venv python (/var/lib/atlas/venv/bin/python, created by
# bootstrap's ensure_atlas_env() — its sanity gate runs THIS exact import before
# any unit points at the venv), not the host python3.
ExecStart=/var/lib/atlas/venv/bin/python -c "import sys; sys.path.insert(0, '/var/lib/atlas/bin'); from atlas.lvm import ThinPool; ThinPool().ensure()"

[Install]
WantedBy=multi-user.target
